# services.d/http-5004/run
#!/usr/bin/with-contenv sh
/usr/local/bin/python /opt/tvhProxy/tvhProxy.py --port 5004 --verbose ${HE_HTTP_VERBOSE:-'0'}